# templates/control_plane/kube-scheduler.service.epp
[Unit]
Description=Kubernetes Scheduler Plugin
Documentation=https://github.com/GoogleCloudPlatform/kubernetes
After=kube-apiserver.service

[Service]
EnvironmentFile=-/etc/kubernetes/common
EnvironmentFile=-/etc/kubernetes/scheduler
User=kube
ExecStart=/opt/hyperkube/bin/kube-scheduler \
            $KUBE_LOGTOSTDERR \
            $KUBE_LOG_LEVEL \
            $KUBE_SCHEDULER_ARGS
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
